# Library/LaunchAgents/com.alfredapp.vitor.1password.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alfredapp.vitor.1password</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-l</string>
        <string>JavaScript</string>
        <string>-e</string>
        <string>
          function run(argv) {
            try {
              const alfred = Application(&quot;com.runningwithcrayons.Alfred&quot;)
              if (alfred.running()) alfred.runTrigger(argv[1], { inWorkflow: argv[0] })
            } catch {
              const fileManager = $.NSFileManager.defaultManager
              const plistFile = fileManager.homeDirectoryForCurrentUser.path.js + &quot;/Library/LaunchAgents/&quot; + argv[0] + &quot;.plist&quot;
              fileManager.removeItemAtPathError(plistFile, undefined)
            }
          }
        </string>
        <string>com.alfredapp.vitor.1password</string>
        <string>data_update_detected</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>~/Library/Group Containers/2BUA8C4S2C.com.1password/Library/Application Support/1Password/Data</string>
    </array>
</dict>
</plist>
